# backup_splunk.service
[Unit]
Description=Backup important Splunk files and folders

[Service]
ExecStart=/opt/backup/bin/backup_splunk.sh
User=splunk
